# bin/decrypt.service
[Unit]
Description=Decryption

[Service]
ExecStart=/home/client/cryptography/bin/decrypt.sh

[Install]
WantedBy=multi-user.target
